# conf/etc/service/statsd2/run
#!/bin/bash

exec /usr/bin/nodejs /opt/statsd/stats.js /opt/statsd/config2.js >> /var/log/statsd2.log 2>&1


